# gw-monitor.service.in
[Unit]
Description=Default Gateway Monitor
Documentation=man:gw-monitor(8)
After=syslog.target

[Service]
Type=simple
EnvironmentFile=-@sysconfdir@/sysconfig/gw-monitor
ExecStart=@sbindir@/gw-monitor $IFNAME
CapabilityBoundingSet=
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
RestrictAddressFamilies=AF_UNIX AF_NETLINK

[Install]
WantedBy=multi-user.target
