# gw-monitor.service.in
[Unit]
Description=Default Gateway Monitor
After=syslog.target

[Service]
Type=simple
EnvironmentFile=-@sysconfdir@/sysconfig/gw-monitor
ExecStart=@sbindir@/gw-monitor $IFNAME

[Install]
WantedBy=multi-user.target
